# reaperd.service
[Unit]
Description=Reaper the auto-deletion daemon

[Service]
ExecStart=/usr/bin/reaperd

[Install]
WantedBy=multi-user.target
